# dot_config/systemd/user/apod.timer
[Unit]
Description=Runs siggy.sh every 5 minutes

[Timer]
OnBootSec=6min
OnUnitActiveSec=6h
Unit=apod.service

[Install]
WantedBy=timers.target
