# dot_config/systemd/user/backup.service
[Unit]
Description=Backup to external drive

[Service]
Type=simple
ExecStart=/home/patrick/.local/bin/homebackup.sh
