# systemd/.config/systemd/user/apod.service
[Unit]
Description=Retrieves APOD from NASA

[Service]
Type=simple
ExecStart=/home/patrick/.local/bin/nasaph.sh

