# systemd/.config/systemd/user/wall.service
[Unit]
Description=Changes wallpaper every 5 minutes

[Service]
Type=simple
ExecStart=/home/patrick/.local/bin/wallpaper.sh

